# reverse-proxy/reverse-proxy.base.service
[Unit]
Description=Reverse Proxy for the Web Services

[Service]
ExecStart=$exec
Restart=always
WorkingDirectory=$working_dir
User=root

[Install]
WantedBy=multi-user.target
