# minecraft/vanilla/1.19.2/minecraft-vanilla.base.service
[Unit]
Description=Minecraft Vanilla Server

[Service]
ExecStart=$exec
WorkingDirectory=$working_dir
KillSignal=SIGINT

[Install]
WantedBy=default.target
